# authwatcher.conf
#this one goes into /etc/init/

description "auth watcher"
start on runlevel [23]
start on started tty1

respawn
respawn limit 10 5
#expect fork
exec /usr/local/authwatcher/authwatcher.py
